# engineos/etc/systemd/system/connman.service.d/disable_dhcp_retries.conf
[Service]
Environment="CONNMAN_DISCOVER_RETRIES=1"
