# engineos/etc/systemd/system/connman.service.d/drop-protection.conf
# Reduce protection because userns doesn't play nicely with overlayfs.
[Service]
ProtectHome=false
ProtectSystem=true
