# engineos/etc/systemd/system/avahi-daemon.service.d/hostname.conf
[Unit]
After=planck-hostname.service
Requires=planck-hostname.service

