# engineos/etc/systemd/system/systemd-udevd.service.d/90-mount-shared.conf
# No history for why Planck configures this when other products don't.
# Is it still needed?
[Service]
MountFlags=shared
