# engineos/usr/lib/systemd/system/az01-power-button.service
[Unit]
Description=AZ01 power button emergency shutdown service

[Service]
Type=simple
ExecStart=/usr/bin/az01-pwrbtn
